# traffic-control.service
[Unit]
Description=Aliyun CDT Traffic Control Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/root/aliyun-traffic-control
ExecStart=/root/aliyun-traffic-control/venv/bin/python /root/aliyun-traffic-control/traffic_control.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
